# scripts/etc/systemd/system/lighttpd.service
# /lib/systemd/system/lighttpd.service
# Managed by QManager — ensures lighttpd uses the correct config.
[Unit]
Description=Lighttpd Daemon
After=network.target opt.mount

[Service]
Type=simple
PIDFile=/opt/var/run/lighttpd.pid
ExecStartPre=/opt/sbin/lighttpd -tt -f /usrdata/qmanager/lighttpd.conf
ExecStart=/opt/sbin/lighttpd -D -f /usrdata/qmanager/lighttpd.conf
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
