# scripts/etc/systemd/system/qmanager-sms-forward.service
# /lib/systemd/system/qmanager-sms-forward.service
[Unit]
Description=QManager SMS Forwarding Daemon
After=network.target qmanager-setup.service

[Service]
Type=simple
ExecStart=/usr/bin/qmanager_sms_forward
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=3600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
