# scripts/etc/systemd/system/qmanager-poller.service
# /lib/systemd/system/qmanager-poller.service
[Unit]
Description=QManager Modem Data Poller
After=qmanager-setup.service qmanager-ping.service
Wants=qmanager-ping.service

[Service]
Type=simple
ExecStartPre=/bin/sh -c '[ -e /dev/smd11 ] || { echo "AT device /dev/smd11 not found"; exit 1; }'
ExecStart=/usr/bin/qmanager_poller
EnvironmentFile=-/etc/qmanager/environment
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=3600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
