# scripts/etc/systemd/system/qmanager-watchcat.service
# /lib/systemd/system/qmanager-watchcat.service
[Unit]
Description=QManager Connection Health Watchdog
After=qmanager-poller.service
Wants=qmanager-poller.service

[Service]
Type=simple
ExecStartPre=/bin/sh -c '[ -f /usr/lib/qmanager/config.sh ] && . /usr/lib/qmanager/config.sh && val=$(qm_config_get watchcat enabled 0) && [ "$val" = "1" ]'
ExecStart=/usr/bin/qmanager_watchcat
EnvironmentFile=-/etc/qmanager/environment
ExecStopPost=/bin/sh -c 'rm -f /tmp/qmanager_watchcat.pid /tmp/qmanager_watchcat.lock /tmp/qmanager_recovery_active'
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=3600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
