# scripts/etc/systemd/system/qmanager-imei-check.service
# /lib/systemd/system/qmanager-imei-check.service
[Unit]
Description=QManager IMEI Rejection Check (One-Shot)
After=qmanager-setup.service

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c '\
    [ -f /etc/qmanager/imei_check_pending ] || exit 1; \
    [ -f /etc/qmanager/imei_backup.json ] || exit 1; \
    enabled=$(jq -r "(.enabled) | if . == null then \"false\" else tostring end" /etc/qmanager/imei_backup.json 2>/dev/null); \
    [ "$enabled" = "true" ] || exit 1'
ExecStart=/usr/bin/qmanager_imei_check
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
